# waiver-video-signage-master/deploy/systemd/signage-update.service
[Unit]
Description=Signage USB Update Service
After=local-fs.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/signage/scripts
ExecStart=/usr/bin/python3 /opt/signage/scripts/usb_update.py

[Install]
WantedBy=multi-user.target
